# dev-db/tarantool-server/files/tarantool.cron
#!/bin/sh

test -x /usr/bin/tarantool_box || exit 0
test -x /usr/bin/tarantool_snapshot_rotate || exit 0

exec /usr/bin/tarantool_snapshot_rotate
